# containers/systemd/salada-app.container
# ~/.config/containers/systemd/salada-app.container
[Container]
ContainerName=salada-app
Image=localhost/salada:latest # <--- REPLACE WITH YOUR ACTUAL IMAGE NAME AND TAG
Pod=salada-main-pod # Link to the same pod as the database

# Environment variable for the database connection string.
# Since both containers are in the same Pod, they share the network namespace.
# So, your Go app can connect to the database via 'localhost' and the DB port.
Environment=DATABASE_URL=postgres://postgres:password@localhost:5432/postgres?sslmode=disable # Adjust based on your Go app's connection string format

# This tells systemd to start the DB container's service before the Go app's service
# Note: This only guarantees start order, not that the DB is fully *ready*.
# Your Go app should ideally have retry logic for DB connections.
Requires=salada-db.container
After=salada-db.container

# Define a basic health check for the Go app (assuming it has a /health endpoint)
HealthcheckIntervalSec=5
HealthcheckTimeoutSec=3
HealthcheckRetries=3
HealthcheckCmd=curl --fail http://localhost:8080/health

# Optional: Resource limits (adjust as needed)
# MemoryLimit=128M
# CPULimit=0.2 # 20% of one CPU core
