# containers/systemd/salada.container
# ~/.config/containers/systemd/salada.container
[Container]
ContainerName=salada-app
IP=192.168.30.10
Image=localhost:5000/salada:latest
Pod=salada # Link to the same pod as the database
Network=salada.network
Requires=salada-db.container
After=salada-db.container
Environment=POSTGRES_HOST=salada-db

# Define a basic health check for the Go app (assuming it has a /health endpoint)
HealthcheckIntervalSec=5
HealthcheckTimeoutSec=3
HealthcheckRetries=3
HealthcheckCmd=curl --fail http://localhost:8080/

# Optional: Resource limits (adjust as needed)
# MemoryLimit=128M
# CPULimit=0.2 # 20% of one CPU core
Secret=postgres-password-kube,type=env,target=POSTGRES_PASSWORD
